# deploy/poiesis.service
# systemd unit for Poiesis. Edit the paths/user, then:
#   sudo cp deploy/poiesis.service /etc/systemd/system/poiesis.service
#   sudo systemctl daemon-reload && sudo systemctl enable --now poiesis
#   journalctl -u poiesis -f
#
# Runs the supervisor, which starts the web app and owns self-mod deploy/rollback.
# systemd only restarts the supervisor itself; the supervisor restarts the app on self-mod.

[Unit]
Description=Poiesis — self-hosted AI
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=matt
WorkingDirectory=/home/matt/poiesis
Environment=HOME=/home/matt
# The Claude Agent SDK shells out to `node` and `claude`; uv runs the app.
# Adjust the node path to your installed version (`dirname $(which node)`).
Environment=PATH=/home/matt/.local/bin:/home/matt/.nvm/versions/node/v25.8.1/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/matt/.local/bin/uv run poiesis start
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
